# etc/systemd/fcch-access-control-door-controller.service
[Unit]
Description=FCCH Access Control Door Controller
Requires=network.target

[Service]
Type=simple
ExecStart=/home/swarren/git_wa/fcch-access-control/bin/door-controller.sh
Restart=always
User=fcchaccess
Group=fcchaccess

[Install]
WantedBy=multi-user.target
